# setup/deployment-poll.cron
* * * * * cd $HOME/deployment && python3 $HOME/manager/src/manager.py
