# setup/deployment-poll.cron
* * * * * cd /home/deployment && python3 /home/manager/manager.py
