# setup/deployment-poll.cron
* * * * * cd /home/ubuntu/deployment && python3 /home/ubuntu/manager/manager.py
